# samples/daemontools/log/run
#!/bin/sh
exec multilog t s500000 n10 /path/to/log/dir
