# contrib/daemontools/log/run
#!/bin/sh
exec multilog s500000 n10 /path/to/log/dir
